# roles/data-container-services-deployer/templates/iplant-data.conf
# Ubuntu upstart entry for iplant-data container

pre-start script
end script

respawn
respawn limit 15 5

start on runlevel [2345]
stop on runlevel [06]

script
/usr/bin/docker run --name {{data_container.container_name}} {{data_container.unc_remote_registry}}/{{data_container.image_name}}:latest
end script
